# .config/systemd/user/ollama.service
[Unit]
Description=Ollama Service

[Service]
ExecStart=%h/bin/ol -s
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
