# .config/systemd/user/photoprism.service
[Unit]
Description=Index Photos

[Service]
Type=simple
Nice=10
ExecStart=/home/sdf/bin/systemd-index-photos
