# .config/systemd/user/restic.service
[Unit]
Description=Restic Backup to Google Storage

[Service]
Type=simple
Nice=10
ExecStart=/home/sdf/bin/rstc
