# .config/systemd/user/status-email-user@.service
[Unit]
Description=Send status email for %i to user

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/systemd-email s@fomichev.me %i
