# .config/systemd/user/sync-email.timer
[Unit]
Description=Email Sync Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=timers.target
