# .config/autostart/watcher.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/watcher
Hidden=false
X-GNOME-Autostart-enabled=true
Name=watcher
Comment=
